# data/systemd/storage-provider-owncloud.service.in
[Unit]
Description=ownCloud storage provider

[Service]
Type=dbus
BusName=com.canonical.StorageFramework.Provider.OwnCloud
ExecStart=@CMAKE_INSTALL_FULL_LIBDIR@/@PROJECT_NAME@/storage-provider-owncloud
