# data/systemd/storage-provider-nextcloud.service.in
[Unit]
Description=Nextcloud storage provider

[Service]
Type=dbus
BusName=com.canonical.StorageFramework.Provider.Nextcloud
ExecStart=@CMAKE_INSTALL_FULL_LIBDIR@/@PROJECT_NAME@/storage-provider-nextcloud
